# mles-webproxy.service.d/override.conf
[Service]
Environment="MLES_KEY=mles-devel-frank"
